# config/crontab
*/10 * * * * /bin/bash -l -c 'cd domains/cloudexchange.org && ruby app/updater.rb >>cron.log 2>&1'
